# deploy/bundle/iot-image-server.service
[Unit]
Description=IoT Image Server - Weather display for e-ink devices
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/bin
EnvironmentFile=%h/bin/env.txt
ExecStart=%h/bin/iot-image-server --lat=${OPEN_WEATHER_LAT} --lon=${OPEN_WEATHER_LON} --open-weather-key=${OPEN_WEATHER_KEY} --stocks-api-key=${TWELVE_DATA_API_KEY}
Restart=on-failure
RestartSec=10s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp

[Install]
WantedBy=default.target
