# deploy/bundle/iot-image-server.service
[Unit]
Description=IoT Image Server - Weather display for e-ink devices
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tboldt
WorkingDirectory=/opt/iot-image-server
EnvironmentFile=/opt/iot-image-server/env.txt
ExecStart=/opt/iot-image-server/iot-image-server --lat=${OPEN_WEATHER_LAT} --lon=${OPEN_WEATHER_LON} --open-weather-key=${OPEN_WEATHER_KEY} --stocks-api-key=${TWELVE_DATA_API_KEY} --stock-symbols=${STOCK_SYMBOLS} --fred-api-key=${FRED_API_KEY}
Restart=on-failure
RestartSec=10s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
